# supervisord.conf
[supervisord]


[program:couchdb-lucene]
directory = ${buildout:directory}/lib/couchdb-lucene-0.10.0-SNAPSHOT
command = java -jar start .jar
